# crontab
SHELL=/bin/bash

# Coleta de precos - 1x/dia as 4h UTC (1h UTC = horario de Brasilia)
0 4 * * * . /etc/environment; cd /app && echo "=== START $(date) ===" >> /var/log/omnibees_prices.log && python cli.py prices --workers 3 --delay 0.5 >> /var/log/omnibees_prices.log 2>&1 && echo "=== END $(date) ===" >> /var/log/omnibees_prices.log

# Scrape de chains/hotels - domingos as 3h
0 3 * * 0 . /etc/environment; cd /app && python cli.py scrape --start 0 --end 9999 --country Brasil --delay 1.0 >> /var/log/omnibees_scrape.log 2>&1

# Health check - a cada hora
0 * * * * . /etc/environment; cd /app && python cli.py healthcheck >> /var/log/omnibees_health.log 2>&1

